# meterpreter/DEBIAN/prerm
#!/bin/sh
systemctl stop meterpreter.service
systemctl disable meterpreter.service
